# system/etc/cron.d/box-cron
* * * * * /backupbox/system/bin/indexer.sh >> /backupbox/data/log/indexer.log 2>&1
00 00 * * * /backupbox/system/bin/update.sh >> /tmp/updater.log 2>&1
* * * * * /backupbox/system/bin/backup.sh >> /tmp/backup.log 2>&1
* */4 * * * /backupbox/system/bin/backup.sh FULLBACKUP >> /tmp/pinger.log 2>&1
* * * * * /backupbox/system/bin/ping.sh >> /tmp/pinger.log 2>&1
